# templates/etc/NetworkManager/dispatcher.d/30-dnscrypt-proxy.sh.j2
#!/bin/sh

if [ "$2" = "up" ]; then
    sleep 5
    systemctl restart dnscrypt-proxy.service
fi