# templates/etc/systemd/system/mirrorlist-update.timer.j2
[Unit]
Description=Mirrorlist Update Timer

[Timer]
Unit=mirrorlist-update.service
OnCalendar=daily
Persistent=True
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
