# templates/etc/systemd/system/sshd.socket.j2

[Unit]
Description=SSH Socket

[Socket]
Accept=true
ListenStream={{ firewall.ports.ssh }}