# templates/etc/systemd/system/startup-custom.service.j2
[Unit]
Description=Commands to run while booting the system

[Service]
ExecStart=/usr/bin/STARTUP_COMMANDS

[Install]
WantedBy=multi-user.target