# templates/etc/systemd/system/docker.service.d/override.conf.j2
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd --data-root {{ docker.storage_root | default('/var/lib/docker') }} -H fd://
   